# root/etc/systemd/system/pacman-mirrors-fasttrack.timer
[Unit]
Description=Run pacman-mirrors --fasttrack daily


[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true


[Install]
WantedBy=timers.target


